# ops/systemd/l6-revoke.timer
[Unit]
Description=Daegis L6 Revoke Evaluator Timer (every 15 min)

[Timer]
OnCalendar=*:0/15
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target