# backups/20251019_1926/halu-self-metrics.service
[Unit]
Description=Halu self-metrics sampler
After=network.target

[Service]
Type=oneshot
User=f
ExecStart=/home/f/daegis/ops/scripts/halu_self_metrics.sh
WorkingDirectory=/home/f

[Install]
WantedBy=multi-user.target
